# s6/services.d/nginx/run
#!/usr/bin/execlineb -P
/opt/nginx/sbin/nginx